# cloudflareufw.timer
[Unit]
Description=Apply cloudflare ip whitelist daily

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timer.target
