# gce-usenet-backup.timer
[Unit]
Description=Backup config to gdrive

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timer.target
